# scripts/backup.cron
# Daily backup for Pocket Money Tracker (ADR-0008).
# Installed to /etc/cron.d/pocket-money-backup by `make install-cron`.

SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

0 3 * * * root /opt/pocket-money/scripts/backup.sh
